# recipes-sota/rvi-sota-client/rvi-sota-client/rvi-sota-client.service
[Unit]
Description=RVI SOTA Client

[Service]
User=root
ExecStart=/usr/bin/run.sh

[Install]
WantedBy=multi-user.target
