# recipes-sota/aktualizr/files/aktualizr-secondary.socket
[Socket]
ListenStream=9030

[Install]
WantedBy=sockets.target